# deploy/tqqq-api.service
[Unit]
Description=TQQQ Trading Bot API Server
After=network.target

[Service]
Type=simple
User=tqqq
Group=tqqq
WorkingDirectory=/opt/tqqq
Environment=PATH=/opt/tqqq/venv/bin:/usr/bin:/bin
ExecStart=/opt/tqqq/venv/bin/python /opt/tqqq/scripts/api_server.py --port 8080
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/tqqq/api.log
StandardError=append:/var/log/tqqq/api-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
